# deploy/targets/run/supervisord.d/nodejs.ini
[program:nodejs]
command=node app.js
directory=/tmp/jenkins_test
stdout_logfile=/tmp/jenkins_test/nodejs.stdout.log
stderr_logfile=/tmp/jenkins_test/nodejs.stderr.log
