# deploy/targets/run/supervisord.d/nodejs.ini
[program:nodejs]
command=node app
directory=/tmp/project/src/app
stdout_logfile=/tmp/project/log/nodejs.stdout.log
stderr_logfile=/tmp/project/log/nodejs.stderr.log
