# deploy/targets/run/supervisord.d/nodejs.ini
[program:nodejs]
command=node app
directory=/mnt/project/
stdout_logfile=/mnt/project/log/nodejs.stdout.log
stderr_logfile=/mnt/project/log/nodejs.stderr.log
